# ASB/ansibleVM/etc/cron.d/ASBTriggerCron
# Cron job that triggers ASBTrigger every minute. Timeout is 30 mins.
# Make sure the symlink to ASBTrigger in /home/ren exists, otherwise this will never get called.

* * * * * root flock -xn /tmp/ASBTrigger.lck /root/ASBTrigger
