# ASB/setup/trigger/etc/cron.d/ASBTriggerCron
# Cron job that triggers ASBTrigger every minute. Timeout is 30 mins.

* * * * * ansible flock -xn /tmp/ASBTrigger.lck /root/ASBTrigger
