# ASB/setup/collector/etc/cron.d/ASBDataCollectorCron
# Cron job that triggers ASBDataCollector script.
# Runs every minute.

*/10 * * * * root /home/ansible/ASBDataCollector
