# ASB/setup/collector/etc/cron.d/ASBDataCollectorCron
# Cron job that triggers ASBDataCollector script.
# Runs every minute.

* * * * * root /home/ansible/ASBDataCollector
