# ASB/cron/ASB.cron
# Timeout after 2 hours if this Ansible has been stuck for more than 2 hours
*/1 * * * * root flock -xn --timeout 7200 /tmp/ASB.lck env GIT_SSL_NO_VERIFY=true ansible-pull --url=http://gerrit.sjc.aristanetworks.com/ardc-config --directory=/root/.ansible/pull ASB/ASB.yml --inventory=localhost,
